# system/arguspi.service
[Unit]
Description=ArgusPI v2 USB Virus Scanner
After=multi-user.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/arguspi
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/opt/arguspi/src
ExecStart=/usr/bin/python3 /opt/arguspi/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=process

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/arguspi/logs /media
ProtectHome=true

[Install]
WantedBy=multi-user.target